# deploy/systemd/mem-mcp-retention-deletion.timer
[Unit]
Description=Run mem-mcp retention_deletion hourly (at :30)
Requires=mem-mcp-retention-deletion.service

[Timer]
OnCalendar=*-*-* *:30:00 UTC
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
