# deploy/systemd/mem-mcp-backup.timer
[Unit]
Description=Run mem-mcp nightly DB backup daily at 02:30 IST
Requires=mem-mcp-backup.service

[Timer]
OnCalendar=*-*-* 21:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
